# configs/supervisor/conf.d/nginx.conf
[program:nginx]
command=/opt/openresty/nginx/sbin/nginx -g 'daemon off;' -c /opt/openresty/nginx/conf/nginx.conf
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0